# services/my-service.service
[Unit]
Description=mytool background service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/my-tool.sh
Restart=on-failure
User=user
Group=user

[Install]
WantedBy=multi-user.target
