# debian/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable my-tool.service
systemctl start my-tool.service



exit 0