# debian/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable my-service.service
systemctl start my-service.service



exit 0